# img/desc/suse-leap-42.1-JeOS/root/etc/systemd/system/getty@tty1.service.d/noclear.conf
[Service]
# ensure tty1 isn't cleared (bnc#804158)
# actually do clear it! DMT
TTYVTDisallocate=yes
